# systemd/dns-updater.timer
[Unit]
Description=Update dynamic DNS every 2 minutes

[Timer]
OnCalendar=*:0/2
Persistent=true

[Install]
WantedBy=timers.target
